# tools/com.stack.claude-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stack.claude-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kayszahidi/bin/refresh-claude-credentials</string>
    </array>

    <!-- Toutes les 2 heures (token expire en ~6-14h) -->
    <key>StartInterval</key>
    <integer>7200</integer>

    <!-- Aussi au démarrage de la session -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/claude-credentials-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-credentials-refresh.log</string>
</dict>
</plist>
